# recollwebui.service
[Unit]
Description=Full text search for personal files
After=network.target
Type=simple

<<<<<<< HEAD
=======
>>>>>>> 5c20953d0b6420b4668adf689b0c30740b9efbb9
# Comentario
# Teste depois de comentario feito em local editado
[Service]
ExecStart=/usr/bin/env python3 /data/AppSys/recollwebui/kmrecollwebui/webui-standalone.py -a %H -p 58143
Restart=on-failure
ExecStop=/bin/kill -SIGINT $MAINPID

[Install]
WantedBy=multi-user.target
